# root/etc/s6-overlay/s6-rc.d/init-permissions-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# initial folders
mkdir -p /data
mkdir -p /config

# permissions
lsiown -R abc:abc \
    /app \
    /config
lsiown abc:abc \
    /data